# deploy/systemd/e_library.service
[Unit]
Description=E-Library FastAPI backend
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/e_library
Environment="ELIBRARY_ROOT=/srv/e_library/library"
Environment="ELIBRARY_LINK_TTL=300"
Environment="ELIBRARY_SIGNING_KEY=CHANGE_ME_TO_STRONG_SECRET"
ExecStart=/opt/e_library/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8085 --workers 2
Restart=always
RestartSec=3

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
